# systemd/measurely-onboarding.service
[Unit]
Description=Measurely onboarding AP fallback
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/measurely-onboarding.sh

[Install]
WantedBy=multi-user.target
